# extra/mcloud-full/debian/mcloud-full.upstart
description "Mcloud full server. Fix common issues with mcloud supporting services."
author "Modera"

start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
    FILE=/var/run/docker.sock
    while [ ! -e $FILE ] ; do
        sleep 1;
    done
    /etc/init.d/dnsmasq start
end script